# assistant/web/conf/supervisor.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
stdout_logfile=/var/log/supervisord.log

[program:nginx]
command=/usr/sbin/nginx

[program:php-fpm]
command=/usr/sbin/php5-fpm
